# amidithru/amidithru-1.0.0/debian/amidithru.service
[Unit]
Description=amidithru service for touchosc2midi
After=sound.target
Wants=sound.target

[Service]
ExecStart=/usr/bin/amidithru touchosc
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target