# mod-ui/mod-ui-1.13.1/debian/mod-update.path
[Unit]
Description=Mod Update Monitor path.

[Path]
PathExists=/tmp/mod-update
Unit=mod-update.service

[Install]
WantedBy=multi-user.target